# assets/autostart/linux-network-speed-indicator.desktop.in
[Desktop Entry]
Type=Application
Version=1.0
Name=Linux Network Speed Indicator
Comment=Live network speed monitor for the Linux top bar and tray
Exec=__EXEC_PATH__
TryExec=/usr/bin/python3
Terminal=false
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=3
OnlyShowIn=GNOME;Unity;
StartupNotify=false
Icon=linux-network-speed-indicator
Categories=Network;
