# assets/autostart/linux-network-speed-indicator.desktop.in
[Desktop Entry]
Type=Application
Version=1.0
Name=Linux Network Speed Indicator
Comment=Show live download and upload speeds in the top bar
Exec=__EXEC_PATH__
TryExec=/usr/bin/python3
Terminal=false
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=3
OnlyShowIn=GNOME;Unity;
StartupNotify=false
Icon=network-transmit-receive-symbolic
Categories=Utility;
